# luci-app-radius-monitor/root/etc/uci-defaults/99_overwrite_build
#!/bin/bash
. /lib/functions.sh

# Cleanup
rm -rf /etc/uci-defaults/99_overwrite_build
exit 0